# luci-app-aliddns/root/etc/init.d/aliddns
#!/bin/sh /etc/rc.common

START=95

restart()
{
	/usr/sbin/aliddns
}

start()
{
	/usr/sbin/aliddns
}

stop()
{
	/usr/sbin/aliddns
}
